# task_project/task_project.service
[Unit]
Description=Task Project Django Application
After=network.target
Wants=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/opt/minicurso-wtads/task_project
Environment=PATH=/opt/minicurso-wtads/task_project/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=DJANGO_SETTINGS_MODULE=task_project.settings
EnvironmentFile=/opt/minicurso-wtads/task_project/.env
ExecStart=/opt/minicurso-wtads/task_project/venv/bin/gunicorn task_project.wsgi:application -c /opt/minicurso-wtads/task_project/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=task_project

[Install]
WantedBy=multi-user.target
